# glances/docker/conf/supervisor.conf
[supervisord]
nodaemon=true
loglevel=info

[program:glances-server]
command=glances -w
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
redirect_stderr=true
