# mosquitto/docker/conf/supervisor/supervisor.conf
[supervisord]
nodaemon=true
loglevel=debug

[program:mosquitto]
command=/usr/sbin/mosquitto -c /config/mosquitto.conf
priority=920
redirect_stderr=true
